# crontab
* * * * * cd $HOME/air-quality-monitor && ./aqm